# app-emacs/telega/files/telega-init
#!/sbin/openrc-run
# Copyright 2008-2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2 or later

# todo: deal with credentials request during the initial run

: ${TELEGA_SERVER:=/usr/bin/telega-server}
USER=${RC_SVCNAME##*.}
: ${TELEGA_SERVER_OPTS:=-l /home/${USER}/.telega/telega-server.log -v 3}
: ${TELEGA_SERVER_TIMEOUT:=30}
PIDFILE_DIR=/var/run/telega-server/${USER}
PIDFILE=${PIDFILE_DIR}/telega-server.pid

description="Start telega-server running in the background"

depend() {
    need localmount
    after bootmisc
}

checkconfig() {
    if [ "${RC_VERSION:-0}" = "0" ]; then
        eerror "This script cannot be used for baselayout-1."
        return 1
    fi

    if [ "${USER}" = "${RC_SVCNAME}" ]; then
	eerror "You have to create an init script for each user:"
	eerror "ln -s telega-server /etc/init.d/telega-server.<user>"
	return 1
    fi

    if ! id -u "${USER}" >/dev/null; then
	eerror "${USER}: No such user"
	return 1
    fi

    checkpath -d --owner 0 --mode 0755 "${PIDFILE_DIR%/*}"
    checkpath -d --owner "${USER}" --mode 0755 "${PIDFILE_DIR}"
}

start() {
    local home
    checkconfig || return 1

    eval home="~${USER}"

    ebegin "Starting telega-server daemon for user ${USER}"
    start-stop-daemon --start \
	--user "${USER}" --pidfile "${PIDFILE}" --chdir "${home}" \
	--exec "${TELEGA_SERVER}" -- ${TELEGA_SERVER_OPTS}
    eend $?
}

stop() {
    ebegin "Stopping telega-server daemon for user ${USER}"
    start-stop-daemon --stop --retry "TERM/${TELEGA_SERVER_TIMEOUT}/KILL/5" \
	--user "${USER}" --pidfile "${PIDFILE}" --exec "${TELEGA_SERVER}"
    eend $?
}
